# cron/crontab
0 * * * * root /app/run_all_scrapers.sh >> /app/logs/cron.log 2>&1
